# init/systemd/kalshi-settlement-poller.service
[Unit]
Description=Kalshi Weather Settlement Poller (30-minute cadence)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/halsted/Documents/python/kalshi_weather
EnvironmentFile=/home/halsted/Documents/python/kalshi_weather/.env
Environment=PYTHONPATH=/home/halsted/Documents/python/kalshi_weather
ExecStart=/home/halsted/Documents/python/kalshi_weather/scripts/run_settlement_poller.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
